# apps/server/deploy/probeflash.service.template
# ProbeFlash later authorized systemd service template.
# Use only after /home/hurricane/probeflash no-sudo verification succeeds and the user authorizes sudo/systemd.
# Copy to /etc/systemd/system/probeflash.service only during S3-SERVER-SYSTEMD-AUTOSTART-VERIFY.

[Unit]
Description=ProbeFlash LAN storage server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hurricane
Group=hurricane
WorkingDirectory=/home/hurricane/probeflash/current/apps/server
EnvironmentFile=/home/hurricane/probeflash/shared/env/probeflash.env
Environment=NODE_ENV=production
ExecStart=/home/hurricane/probeflash/runtime/node/bin/node src/server.mjs
Restart=always
RestartSec=3s

# Logs: journald still records unit output; these append files make handoff/debug copy easier.
StandardOutput=append:/home/hurricane/probeflash/shared/logs/probeflash.log
StandardError=append:/home/hurricane/probeflash/shared/logs/probeflash.err.log

# Hardening: keep the app read-only except explicit shared writable paths.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/hurricane/probeflash/shared/data /home/hurricane/probeflash/shared/logs /home/hurricane/probeflash/shared/env
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
